# btc-dashboard.service
[Unit]
Description=BTC Dashboard Streamlit Application
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/btc-dashboard
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="DASHBOARD_PORT=8502"
ExecStart=/usr/bin/python3 -m streamlit run app.py --server.port=${DASHBOARD_PORT} --server.address=127.0.0.1 --server.headless=true --server.enableCORS=false --server.enableXsrfProtection=false
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=yes
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/btc-dashboard

[Install]
WantedBy=multi-user.target
